# docs/install/com.turtle.river.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.turtle.river</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/turtle/turtleos/venv/bin/python3</string>
        <string>-u</string>
        <string>/Users/turtle/turtleos/river_bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/turtle/turtleos</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>DOTENV_PATH</key>
        <string>/Users/turtle/turtleos/.env</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/turtle/turtleos/logs/river.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/turtle/turtleos/logs/river.log</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
